# ops/com.connor.recipes-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.connor.recipes-daily</string>

  <!-- zsh -lc reads ~/.zprofile (homebrew: node, gh) but NOT ~/.zshrc, which is
       interactive-only and is where ~/.local/bin (claude) is added. Prepend it
       explicitly or the job dies with "claude: command not found".
       Note: && must be escaped as &amp;&amp; to keep this valid XML. -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/zsh</string>
    <string>-lc</string>
    <string>export PATH="$HOME/.local/bin:$PATH" &amp;&amp; cd "$HOME/recipes-that-dont-exist" &amp;&amp; git pull --rebase --quiet &amp;&amp; claude -p "/daily-recipe" --dangerously-skip-permissions</string>
  </array>

  <!-- 18:00 local is past 00:00 UTC year-round for US Pacific
       (PDT 18:00 = 01:00 UTC; PST 18:00 = 02:00 UTC). Adjust if the
       mini lives in another timezone. -->
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>18</integer>
    <key>Minute</key>
    <integer>5</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/tmp/recipes-daily.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/recipes-daily.log</string>
</dict>
</plist>
